# systemd/telework-task.service
[Unit]
Description=telework-tasks
After=network.target

[Service]
Type=simple
ExecStart=java NetworkReceiver
WorkingDirectory=/home/telework/tasks
User=telework
Group=telework
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
